# timelimit.conf
start on runlevel [2345]
stop on runlevel [016]
respawn
script
  exec /home/ags/projects/linux-timeacc/timelimit.sh
end script

